# scripts/lemmatizer.service.example
[Unit]
Description=Pymorphy2 Lemmatizer Service
After=network.target

[Service]
Type=simple
# Replace YOUR_USER with your username, or use setup_lemmatizer.sh
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/patchouli/bot/scripts
ExecStart=/home/YOUR_USER/patchouli/bot/scripts/venv/bin/python3 /home/YOUR_USER/patchouli/bot/scripts/lemmatizer_service.py --socket /tmp/lemmatizer.sock
Restart=on-failure
RestartSec=5

# Cleanup socket on stop
ExecStopPost=/bin/rm -f /tmp/lemmatizer.sock

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=false
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
